# fstbot.service
[Unit]
Description=Slack fstbot service

[Service]
Type=simple
ExecStart=/bin/bash fstbot.sh

[Install]
WantedBy=multi-user.target
